# services.d/transmission/run
#!/usr/bin/execlineb -P

with-contenv
s6-envuidgid transmission

multisubstitute
{
  importas -i config TRANSMISSION_CONFIG
}

transmission-daemon --foreground --config-dir ${config}
